# app/front/manaul.service
<<<<<<< HEAD
[Unit]
Description=Token Manual Entry Server

[Service]
ExecStart=/home/fullsend/cryptovesting/app/front/start_server.sh

[Install]
=======
[Unit]
Description=Token Manual Entry Server

[Service]
ExecStart=/home/fullsend/cryptovesting/app/front/start_server.sh

[Install]
>>>>>>> 1e3863f015956128c6e5f2af8081280410a04cc7
WantedBy=multi-user.target